# com.ecosystem.notion-poller.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ecosystem.notion-poller</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/mdmac/Documents/ecosystem-mcp-server/.venv/bin/python</string>
        <string>-m</string>
        <string>ecosystem_mcp_server.notion_control</string>
        <string>--interval</string>
        <string>60</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/mdmac/Documents/ecosystem-mcp-server</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/mdmac/Library/Logs/ecosystem-notion-poller.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/mdmac/Library/Logs/ecosystem-notion-poller.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
